# infra/systemd/merchants.service
[Unit]
Description=Merchants Host (Saving payment platform)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/root/app/merchants
Restart=on-failure
RestartSec=5

Environment=MERCHANTS_DB=/var/lib/docker/volumes/app_merchants_data/_data/merchants.db
Environment=MERCHANTS_ADDR=:8090
Environment=WIRE_ADMIN_URL=http://localhost:7475
Environment=WIRE_M2M_TOKEN=REPLACE_WITH_SECRET

WorkingDirectory=/root/app
StandardOutput=append:/var/log/merchants.log
StandardError=append:/var/log/merchants.log

[Install]
WantedBy=multi-user.target
